# mkosi-extra/etc/systemd/system/arexibo-update.timer
[Unit]
Description=Check for arexibo updates

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
