# mkosi-extra/etc/systemd/system/arexibo-update.service
[Unit]
Description=Update arexibo package

[Service]
Type=oneshot
ExecStart=/usr/bin/dnf update -y arexibo
